# xplane12/host/systemd/xplane12-data-api.service
[Unit]
Description=X-Plane 12 local snapshot API
After=network-online.target xplane12-autoflight.service
Wants=network-online.target
Requires=xplane12-autoflight.service

[Service]
Type=simple
User=tang
EnvironmentFile=-/home/tang/xplane12.env
ExecStart=/home/tang/xplane12_data_api.sh
Restart=always
RestartSec=5
Environment=HOME=/home/tang
WorkingDirectory=/home/tang

[Install]
WantedBy=multi-user.target
