# xplane12/host/systemd/xplane12-autoflight.service
[Unit]
Description=X-Plane 12 endless-flight runtime
After=network-online.target xplane12-simulator.service
Wants=network-online.target
Requires=xplane12-simulator.service

[Service]
Type=simple
User=your-user
EnvironmentFile=-/home/your-user/xplane12.env
ExecStart=/home/your-user/xplane12_autoflight.sh
Restart=always
RestartSec=5
Environment=HOME=/home/your-user
WorkingDirectory=/home/your-user

[Install]
WantedBy=multi-user.target
