# xplane12/host/systemd/xplane12-simulator.service
[Unit]
Description=X-Plane 12 simulator
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=your-user
EnvironmentFile=-/home/your-user/xplane12.env
ExecStart=/home/your-user/xplane12_launch.sh
Restart=always
RestartSec=20
Environment=HOME=/home/your-user
WorkingDirectory=/home/your-user

[Install]
WantedBy=multi-user.target
