# xplane12/host/systemd/xplane-49013-tunnel.service
[Unit]
Description=X-Plane TCP reverse tunnel on port 49013
After=network-online.target xplane12-simulator.service
Wants=network-online.target
Requires=xplane12-simulator.service

[Service]
Type=simple
User=your-user
EnvironmentFile=-/home/your-user/xplane12.env
ExecStart=/home/your-user/tunnel_xplane_49013.sh
Restart=always
RestartSec=5
Environment=HOME=/home/your-user
WorkingDirectory=/home/your-user

[Install]
WantedBy=multi-user.target
